# feeds/luci/applications/luci-goagent/root/etc/uci-defaults/luci-goagent
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@goagent[-1]
	add ucitrack goagent
	set ucitrack.@goagent[-1].init=goagent
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
